# system_files/shared/etc/profile.d/common.sh
alias open="xdg-open &>/dev/null"
alias neofetch=fastfetch

_eval_if_available starship init bash
