# system_files/shared/etc/profile.d/update.sh
alias rollback="sudo bootc rollback"
alias update="spacefin-cli run update"
